# local-routing/systemd/health-router.timer
[Unit]
Description=Run health-router periodically

[Timer]
OnBootSec=30
OnUnitActiveSec=15

[Install]
WantedBy=timers.target